# debian/dovecot-pop3d.postinst
#!/bin/sh

set -e

#perl -pi.bak -e 'if (/^\s*protocols =/i) { s/$/ pop3 pop3s/ unless /pop3/; };'\
#   /etc/dovecot/dovecot.conf
if [ -x "/etc/init.d/dovecot" ]; then
	if [ -x /usr/sbin/invoke-rc.d ]; then
		invoke-rc.d dovecot start
	else
		/etc/init.d/dovecot start
	fi
fi

#DEBHELPER#
