# debian/dovecot-pop3d.prerm
#!/bin/sh
set -e

if [ -d /run/systemd/system ]; then
	deb-systemd-invoke stop dovecot.socket || true
fi

if [ -x "/etc/init.d/dovecot" ]; then
	if [ -x invoke-rc.d ] ; then
		invoke-rc.d dovecot stop
	else
		/etc/init.d/dovecot stop
	fi
fi

#DEBHELPER#
