# debian/dovecot-pop3d.postinst
#!/bin/sh

set -e

if [ "$1" = "configure" -a -z "$2" ]; then
	# Add the pop3 option to the protocols line on first install
	sed -i -e '/^[[:space:]]*protocols *=/ {s/[[:space:]]*none//; /pop3/!s/[[:space:]]*$/ pop3/}' \
		/etc/dovecot/protocols.conf
fi

if [ -x "/etc/init.d/dovecot" ]; then
	if [ -x /usr/sbin/invoke-rc.d ]; then
		invoke-rc.d dovecot restart
	else
		/etc/init.d/dovecot restart
	fi
fi

#DEBHELPER#
