# debian/dovecot-lmtpd.prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/dovecot" ]; then
	if [ -x invoke-rc.d ] ; then
		invoke-rc.d dovecot stop
	else
		/etc/init.d/dovecot stop
	fi
fi

#DEBHELPER#
